# kicker.conf
start on file FILE=/home/USER/DIRECTORY/success.txt EVENT=create
setuid USERNAME
console log
chdir /home/USER/DIRECTORY/
task
	exec /home/USER/DIRECTORY/master.sh 
